# Yocto/meta-infotainment/recipes-kernel/load-modules/files/load-modules.service
[Unit]
Description=Load kernel drivers

[Service]
Type=oneshot
ExecStart=/usr/bin/load-modules.sh

[Install]
WantedBy=multi-user.target
